# systemd.timer.in
[Unit]
Description=Runs the day/night theme switcher every 10 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
Unit=SERVICE_NAME.service

[Install]
WantedBy=default.target
